# system_files/overrides/usr/lib/systemd/system/return-to-gamemode.service
[Unit]
Description=Return to HOSTNO4 Mode

[Service]
Type=oneshot
ExecStart=/usr/bin/return-to-gamemode
